# debian/pihsm-display.cron.hourly
#!/bin/sh -e

systemctl restart pihsm-display.service
