# debian/pihsm-server.service
[Unit]
Description=pihsm server

[Service]
User=pihsm-server
ExecStart=/usr/bin/pihsm-server
Restart=always

[Install]
WantedBy=multi-user.target

