# debian/pihsm-server.cron.hourly
#!/bin/sh -e

systemctl restart pihsm-server.service
